# root/etc/cron.daily/phonebook
#!/bin/bash

/usr/share/phonebooks/phonebook 2&> /dev/null
