# root/etc/cron.daily/phonebook
#!/bin/bash

/usr/share/phonebooks/phonebook &> /dev/null
/sbin/restart slapd &> /dev/null
