# root/etc/cron.d/phonebook
30 13,23 * * * root /usr/share/phonebooks/phonebook
